# systemd_units/homeguarantee.service
[Unit]
Description=Automatically restarts HomeAssistant
Requires=homeguarantee.timer

[Service]
ExecStart=docker containers start homeassistant

