# systemd_units/homeguarantee.timer
[Unit]
Description=Ensure HomeAssistant

[Timer]
Unit=kieltimer.service
OnCalendar=*-*-* 03:30:30

[Install]
WantedBy=timers.target

