# systemd_units/kieltimer.timer
[Unit]
Description=Run something

[Timer]
Unit=kieltimer.service
OnCalendar=*-*-* *:00:30

[Install]
WantedBy=timers.target

