# systemd_units/kiel.target
[Unit]
Description=Kiel units
Requires=kieltimer.service
Requires=kielserver.service
Requires=kielfetch.service

