# systemd_units/kieltimer.service
[Unit]
Description=Runs something
Wants=kieltimer.timer

[Service]
ExecStart=/home/pi/kieltimer.sh

[Install]
WantedBy=default.target

