# systemd_units/kielserver.service
[Unit]
Description=Kiel server Daemon

[Service]
ExecStart=/usr/local/bin/kielserver

[Install]
WantedBy=kiel.target

