# systemd_units/kielfetch.service
[Unit]
Description=Fetches Kiel data

[Service]
ExecStart=/usr/local/bin/kiel fetch --enact

[Install]
WantedBy=multi-user.target

