# systemd_units/kiel.target
[Unit]
Description=Kiel units
Requires=kieltimer.timer
Requires=kielserver.service
# Requires=kielfetch.timer

[Install]
WantedBy=multi-user.target

