# systemd_units/kieltimer.service
[Unit]
Description=Runs something
Requires=kieltimer.timer

[Service]
ExecStart=/usr/local/bin/kiel --hour

[Install]
WantedBy=kiel.target

