# systemd_units/kieltimer.timer
[Unit]
Description=Run something

[Timer]
Unit=kieltimer.service
OnCalendar=*-*-* *:00:30
AccuracySec=1s

[Install]
WantedBy=timers.target

