# systemd_units/kielserver.service
[Unit]
Description=Kiel server Daemon

[Service]
ExecStart=python /usr/local/bin/kielserver.py

[Install]
WantedBy=kiel.target

