# systemd_units/kiel.target
[Unit]
Description=Kiel units
Requires=kieltimer.timer
Requires=kielserver.timer
# Requires=kielfetch.service

[Install]
WantedBy=multi-user.target

